# update_cloudflare_ips.cron
# Run update every 12 hours
0 */12 * * * root /usr/local/bin/update_cloudflare_ips.sh
